# Tasks/Task_1.3/uptime-logger.service
[Unit]
Description=Uptime Logger Service
After=network.target

[Service]
ExecStart=/home/ubuntu/uptime-logger.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
